# systemd/hotdeal-discovery.service
[Unit]
Description=Hotdeal Bot Discovery Job
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/hyeonbin/hotdeal_bot
ExecStart=/home/hyeonbin/hotdeal_bot/.venv/bin/python /home/hyeonbin/hotdeal_bot/src/jobs/run_discovery.py
Restart=no

